# deploy/iroh-relay.service
[Unit]
Description=iroh relay server (self-hosted, server.viroh.net)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Path from `cargo install iroh-relay --features server` run as the ubuntu user.
# (Runs as root via User=root below so it can bind :443 and read the TLS key.)
ExecStart=/home/ubuntu/.cargo/bin/iroh-relay --config-path /home/ubuntu/viroh/deploy/iroh-relay.config.toml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
